# gunicorn.conf

[program:gunicorn]

directory = /home/ubuntu/django-pd/django_ref/
command = /home/ubuntu/django-pd/env/bin/gunicorn --workers=2 --bind unix:/home/ubuntu/django-pd/django_ref/app.sock django_ref.wsgi:application --timeout 400
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log
